# systemd/kb-bot.service
[Unit]
Description=Signal KB ingestion bot
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/signal-kb
EnvironmentFile=/root/signal-kb/.env
ExecStart=/root/signal-kb/.venv/bin/python -m pipeline.bot
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
